# modules/mark2/templates/service.sh.erb
#!/bin/bash
# /etc/init.d/<%= @title %>

### BEGIN INIT INFO
## Provides: minecraft
## Required-Start: $local_fs $remote_fs
## Required-Stop: $local_fs $remote_fs
## Should-Start: $network
## Should-Stop: $network
## Default-Start: 2 3 4 5
## Default-Stop: 0 1 6
## Short-Description: <%= @title %> Minecraft server
## Description: Init script for the <%= @title %> Minecraft server
#### END INIT INFO

# Created by Nullreff

# User that should run the server
MC_USER="<%= @user %>"
# The user running this script
CURRENT_USER=`whoami`
# BungeeCord home direcory
MC_PATH="<%= @home %>"
# Server server
MC_NAME=`basename $MC_PATH`

run_as() {
    if [ $CURRENT_USER == $MC_USER ]; then
        bash -c "$1"
    elif command -v sudo >/dev/null 2>&1; then
        sudo su - $MC_USER -c "$1"
    else
        su - $MC_USER -c "$1"
    fi
}

case "$1" in
    start)
        run_as "
        <%- if @port != '25565' -%>
            # Make sure the server is running on the right port
            sed -i 's/server-port=[0-9]*/server-port=<%= @port %>/' $MC_PATH/server.properties

        <%- end -%>
        <%- if @bungeeCord -%>
            # Enable bungeecord settings
            sed -i 's/online-mode=true/online-mode=false/' $MC_PATH/server.properties
            sed -i 's/connection-throttle: [0-9]*/connection-throttle: -1/' $MC_PATH/server.properties

        <%- end -%>
            mark2 start $MC_PATH
        "
        ;;
    stop)
        run_as "mark2 stop -n $MC_NAME"
        ;;
    status)
        if mark2 list | grep $MC_NAME > /dev/null; then
            echo "<%= @title %> is running"
            exit 0
        else
            echo "<%= @title %> is not running"
            exit 1
        fi
        ;;
    *)
        echo "Usage: /etc/init.d/<%= @title %> {start|stop|status}"
        exit 1
        ;;
esac

exit 0
